# etc/supervisord.d/graylog2-web-interface.ini
[program:graylog2-web-interface]
command = /opt/graylog2-web-interface/bin/graylog2-web-interface
startsecs=20
redirect_stderr=true
stdout_logfile = /var/log/supervisor/%(program_name)s.log
autorestart = true
