# etc/supervisord.d/graylog2-server.ini
[program:graylog2-server]
command = /bin/java -jar /opt/graylog2-server/graylog2-server.jar
startsecs=30
redirect_stderr=true
stdout_logfile = /var/log/supervisor/%(program_name)s.log
autorestart = true
