# etc/supervisord.d/mongodb.ini
[program:mongodb]
command=/bin/mongod --dbpath /data/db --bind_ip '127.0.0.1' --journal  --smallfiles
priority = 10
redirect_stderr=true
stdout_logfile = /var/log/supervisor/%(program_name)s.log
autorestart = true
